# build-docker-php-5.6/config/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/bin/php-fpm5

[program:nginx]
command=/usr/sbin/nginx
